# skel/etc/profile.d/default.sh
export HOSTNAME=`hostname`
export PATH="$PATH:/picar/bin:/picar/sbin"
export PS1="\033[1;31m$HOSTNAME\033[0m:\w> "
